# bin/cfb_daemon_task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>CFB Ticket Tracker price scraping and modeling daemon</Description>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <RestartOnFailure>
      <Interval>PT2M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
    <StartWhenAvailable>true</StartWhenAvailable>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\randi\OneDrive - Clemson University\Desktop\Randi Folders\Personal\Football Ticket Tracker\cfb-ticket-tracker\.venv_win\Scripts\python.exe</Command>
      <Arguments>-m cfb_tix run --no-gui</Arguments>
      <WorkingDirectory>C:\Users\randi\OneDrive - Clemson University\Desktop\Randi Folders\Personal\Football Ticket Tracker\cfb-ticket-tracker</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
